# executor.service
[Unit]
Description=t3rn Executor Node
After=network.target

[Service]
User=tokio
WorkingDirectory=/opt/t3rn-executor
ExecStart=/opt/t3rn-executor/executor --executor --network sepolia
Restart=always
RestartSec=10
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
